# systemd/1ai-engage-watchdog.service
[Unit]
Description=1ai-engage Watchdog
After=1ai-engage-mcp.service 1ai-engage-ui.service

[Service]
Type=simple
User=openclaw
ExecStart=/home/openclaw/.openclaw/workspace/1ai-engage/scripts/watchdog.sh
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
