# systemd/1ai-engage-autonomous.service
[Unit]
Description=1ai-engage Autonomous Outreach Loop
After=network.target 1ai-engage-mcp.service

[Service]
Type=simple
User=openclaw
WorkingDirectory=/home/openclaw/.openclaw/workspace/1ai-engage
ExecStart=/home/linuxbrew/.linuxbrew/bin/python3 /home/openclaw/.openclaw/workspace/1ai-engage/scripts/autonomous_loop.py
Restart=always
RestartSec=10
StandardOutput=append:/home/openclaw/.openclaw/workspace/1ai-engage/logs/autonomous.log
StandardError=append:/home/openclaw/.openclaw/workspace/1ai-engage/logs/autonomous.log

[Install]
WantedBy=multi-user.target
